# tether/tether-iptables.service
[Unit]
Description=tether-iptables

[Service]
ExecStart=/usr/local/sbin/tether

[Install]
WantedBy=multi-user.target
